# supervisord.conf
[supervisord]
nodaemon=true

[program:django_kafka]
command=python -u manage.py consume_kafka
autostart=true
autorestart=true
logfile=%(here)s
